# package/base-files/files/etc/uci-defaults/disable-adslcheck-default
#!/bin/sh
/etc/init.d/adsl_check stop
/etc/init.d/adsl_check disable
exit 0
