# package/base-files/files/etc/uci-defaults/fix-dropbear-perms
#!/bin/sh
chmod 700 /etc/dropbear
chmod 600 /etc/dropbear/*
exit 0
