# package/base-files/files/etc/uci-defaults/fix-dnsmasq-rebind-protection
uci set dhcp.@dnsmasq[0].rebind_protection=0
uci commit
